# templates/hashicorp/vault.service.erb
[Unit]
Description=Unit that keeps vault server up and running

[Service]
WorkingDirectory=<%= @homeDir %>
ExecStart=/usr/local/bin/vault server -config <%= @vaultConfig %>
Restart=always

[Install]
WantedBy=multi-user.target
